# camera.service
[Unit]
Description=Prusa Camera Timelapse Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/prusa-camera
ExecStart=/usr/bin/python3 /home/pi/prusa-camera/camera_service.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

